# etc/supervisord.d/diamond.ini
[program:diamond]
command=/bin/diamond -f 
redirect_stderr=true
stdout_syslog=true
