# etc/supervisord.d/elasticsearch.ini
[program:elasticsearch]
environment=ES_HEAP_SIZE="2g",ES_MIN_MEM="2g",ES_MAX_MEM="2g"
command=/usr/share/elasticsearch/bin/elasticsearch -f -p /var/run/elasticsearch/elasticsearch.pid \
          -Des.default.path.home=/usr/share/elasticsearch -Des.default.path.logs=/var/log/elasticsearch \
          -Des.default.path.data=/var/lib/elasticsearch -Des.default.path.work=/tmp/elasticsearch -Des.default.path.conf=/etc/elasticsearch
redirect_stderr=true
stdout_syslog=true
autorestart = true
priority = 10
