# etc/supervisord.d/kibana.ini
[program:kibana]
command = /opt/kibana/bin/kibana
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
