# etc/supervisord.d/logstash_watchdog.ini
[program:logstash_watchdog]
command=/root/bin/watch_logstash.sh
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
#stdout_syslog=true
