# etc/supervisord.d/nginx.ini
[program:nginx]
command = /opt/qnib/bin/start_nginx.sh
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
