# etc/supervisord.d/nginx.ini
[program:nginx]
command = /usr/sbin/nginx
redirect_stderr=true
stdout_syslog=true
