# etc/supervisord.d/setup.ini
[program:setup]
command=/usr/local/bin/qnib-setup.py
priority=1
startretries=5
startsecs=5
redirect_stderr=true
stdout_syslog=false
stdout_logfile = /var/log/supervisor/%(program_name)s.log
