# etc/supervisord.d/setup.ini
[program:setup]
command=/root/bin/setup.sh
priority=1
startretries=0
startsecs=0
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

