# etc/supervisord.d/sshd.ini
[program:sshd]
command=/sbin/sshd -D
redirect_stderr=true
stdout_syslog=true
