# etc/supervisord.d/statsd.ini
[program:statsd]
command=/usr/bin/statsd /etc/statsd/config.js
redirect_stderr=true
stdout_syslog=true
