# etc/supervisord.d/syslog-ng.ini
[program:syslog-ng]
command=/usr/sbin/syslog-ng --foreground
#stdout_syslog=true
autostart=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stdout_capture_maxbytes=1MB
